# docker/pulsar-dashboard/pulsar-dashboard.conf

[program:pulsar-ui-postgres]
command = /usr/lib/postgresql/9.6/bin/postgres -D /etc/postgresql/9.6/main
user = postgres

[program:pulsar-ui]
directory = /pulsar/django
command = /pulsar/django/manage.py runserver 0.0.0.0:8000

[program:pulsar-stats-collector]
command = /pulsar/django/collector.sh
